# ansible/roles/runtime/templates/sask.service.j2
[Unit]
Description=sask gunicorn service
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_root }}
EnvironmentFile={{ etc_dir }}/environment
# No --access-logfile/--error-logfile (SPEC-032): gunicorn's stdout/stderr
# under systemd is captured by journald as-is. The app's own structured
# JSON (src/sask/logsetup.py) is the request-level record; a separate
# gunicorn access log would duplicate it and mix plain-text lines into
# journald's JSON-per-line stream. gunicorn's own stderr (worker
# boot/reload/error notices) still flows to journald via systemd's default
# stream capture.
# --no-control-socket: gunicorn 26+'s control socket defaults under the
# user's home dir, which ProtectHome=true blocks (surfaced as a harmless
# but noisy "Control server error: Permission denied" on every start once
# --error-logfile stopped swallowing it into an unread file); this service
# doesn't use the control interface, so disable it rather than carve out
# a ProtectHome exception for an unused feature.
ExecStart={{ venv_dir }}/bin/gunicorn {{ wsgi_module }} \
    --bind {{ gunicorn_bind }} \
    --workers {{ gunicorn_workers }} \
    --timeout {{ gunicorn_timeout }} \
    --max-requests {{ gunicorn_max_requests }} \
    --max-requests-jitter {{ gunicorn_max_requests_jitter }} \
    --no-control-socket
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
